# systemd/user/theme-sun-switch.service
[Unit]
Description=Cambia tema claro/oscuro según el sol

[Service]
ExecStart=%h/.config/scripts/theme-sun-switch.sh
